# dotfiles/mopidy/mopidy.service
[Unit]
Description=Mopidy music server
After=dbus.service
After=sound.target

[Service]
ExecStart=mopidy

[Install]
WantedBy=default.target
