# roles/bat/files/profile.d/10-bat.sh
alias -g cat="bat --paging=never"
alias -g less="bat --paging=always"
