# roles/fd/files/profile.d/10-fd.sh
# shellcheck disable=SC1091

[[ -z "$LS_COLORS" ]] && source /usr/share/LS_COLORS/dircolors.sh

alias find="fd"
