# roles/bob/files/profile.d/50-bob-path.sh
export PATH="$PATH:$XDG_DATA_HOME/bob-install"
