# roles/unclutter/files/systemd/unclutter.service
[Unit]
Description=Hide mouse cursors on inactivity
After=xorg-server.service

[Service]
ExecStart=unclutter

[Install]
WantedBy=xorg-server.service
