# dotfiles/unclutter/unclutter.service
[Unit]
Description=Hide mouse cursors on inactivity

[Service]
ExecStart=unclutter

[Install]
WantedBy=default.target
