# dotfiles/unclutter/unclutter.service
[Unit]
Description=Hide mouse cursors on inactivity

[Service]
ExecStart=echo Display $DISPLAY

[Install]
WantedBy=default.target
