# roles/polybar/files/systemd/suspend-polybar.service
[Unit]
Before=suspend.target

[Service]
Type=simple
ExecStart=/bin/systemctl stop polybar.service

[Install]
WantedBy=suspend.target
