# roles/zathura/files/profile.d/90-alias-zathura.sh
alias zathura="tabbed -dcn 'Zathura Tabbed' -p -1 | xargs -I XID zathura --reparent=XID"
