# session/autostart/compton.desktop
[Desktop Entry]
Name=Compton
GenericName=A compositor for X11
Exec=compton --daemon
Terminal=false
Type=Application
