# session/autostart/compton.desktop
[Desktop Entry]
Name=Compton
GenericName=A compositor for X11
Exec=compton -b -d:0
Terminal=false
Type=Application
