# session/autostart/dunst.desktop
[Desktop Entry]
Name=Dunst
GenericName=A customizable and lightweight notification-daemon
Exec=dunst
Terminal=false
Type=Application
