# session/autostart/mopidy.desktop
[Desktop Entry]
Name=Mopidy
GenericName=An extenable music player daeomon
Exec=mopidy
Terminal=false
Type=Application
