# session/autostart/powerstatus10k.desktop
[Desktop Entry]
Name=Powerstatus10k
GenericName=Powerful statusbar
Exec=powerstatus10k
Terminal=false
Type=Application
