# session/autostart/powerstatus10k.desktop
[Desktop Entry]
Name=Powerstatus10k
GenericName=Powerful statusbar
Exec=/home/thore/Tools/powerstatus10k/powerstatus10k.sh
Terminal=false
Type=Application
