# session/autostart/redshift.desktop
[Desktop Entry]
Name=Redshift
GenericName=Display color temperature adjustment
Exec=/usr/bin/redshift
Terminal=false
Type=Application
