# session/autostart/wallpaper.desktop
[Desktop Entry]
Name=Wallpaper
GenericName=Set the X11 wallaper
Exec=feh --bg-scale /home/thore/Bilder/Wallpapers/Wallpaper.png
Terminal=false
Type=Application
