# roles/X11/files/xorg-server.service
[Unit]
Description=Xorg server

[Service]
Type=oneshot
ExecStart=echo 'This is just a dump unit used to artificially syncing dependent services'
