# roles/xbindkeys/files/systemd/xbindkeys.service
[Unit]
Description=Key bindings for X11

[Service]
ExecStart=xbindkeys --nodaemon --file $XBINDKEYS_CONFIG_FILE

[Install]
WantedBy=default.target
