# roles/powerlevel10k/files/profile.d/89-powerlevel10k-load.zsh
# All profiles that want to add segments on top of the default configuration
# should order before this one here. Preferably they use the 8x indexes. The
# segment order depends on the indexes. Show on command segments are preferred.

source /usr/share/zsh-theme-powerlevel10k/powerlevel10k.zsh-theme
