# old_to_migrate/dotfiles/autostart/orgatimer.desktop
[Desktop Entry]
Name=Orgatime
GenericName=Timer to update/sync all organizer tools
Exec=orgatimer
Terminal=false
Type=Application
