# dotfiles/autostart/compton.desktop
[Desktop Entry]
Name=Compton
GenericName=A compositor for X11
Exec=compton --daemon --corner-radius 10
Terminal=false
Type=Application
