# old_to_migrate/dotfiles/picom/picom.service
[Unit]
Description=Window compositor for X11

[Service]
ExecStart=picom

[Install]
WantedBy=default.target
